# ramadan-scheduler.service.example
# systemd unit for Ramadan scheduler (Ubuntu / Linux)
# Run at sunrise and sundown for the duration of Ramadan.
#
# 1. Copy to systemd user or system dir:
#    cp ramadan-scheduler.service.example ~/.config/systemd/user/ramadan-scheduler.service
#    # or for system-wide: sudo cp ramadan-scheduler.service.example /etc/systemd/system/
#
# 2. Edit the file: set your project path and Python path (or venv).
#
# 3. If using a venv, set ExecStart to:
#    /path/to/RamadanOeroen/.venv/bin/python /path/to/RamadanOeroen/ramadan_scheduler.py
#
# 4. Optional: create an env file for Discord/location:
#    /path/to/RamadanOeroen/ramadan.env with RAMADAN_LATITUDE, RAMADAN_LONGITUDE, RAMADAN_TIMEZONE, DISCORD_WEBHOOK_URL
#
# 5. Enable and start (user):
#    systemctl --user daemon-reload
#    systemctl --user enable ramadan-scheduler
#    systemctl --user start ramadan-scheduler
#
#    Or system-wide:
#    sudo systemctl daemon-reload
#    sudo systemctl enable ramadan-scheduler
#    sudo systemctl start ramadan-scheduler
#
# 6. Check status and logs:
#    systemctl --user status ramadan-scheduler
#    journalctl --user -u ramadan-scheduler -f

[Unit]
Description=Ramadan video scheduler (sunrise and sundown)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/path/to/RamadanOeroen
ExecStart=/usr/bin/python3 /path/to/RamadanOeroen/ramadan_scheduler.py
Restart=on-failure
RestartSec=60

# Optional: load env vars (Discord webhook, location)
# EnvironmentFile=/path/to/RamadanOeroen/ramadan.env

# Optional: if using a virtualenv
# ExecStart=/path/to/RamadanOeroen/.venv/bin/python /path/to/RamadanOeroen/ramadan_scheduler.py

[Install]
WantedBy=default.target
